# services/hermes-voice-agent.service
[Unit]
Description=Hermes Voice Agent - Persistent voice presence
After=network.target sound.target llama-voice.service
Requires=llama-voice.service

[Service]
Type=simple
User=sysadmin
WorkingDirectory=/home/sysadmin
ExecStartPre=/usr/local/bin/alienware-audio-fix.sh
EnvironmentFile=/home/sysadmin/.hermes/.env
ExecStart=/home/sysadmin/.hermes/venv/bin/python3 /home/sysadmin/hermes-voice-agent.py
Restart=no
StandardOutput=append:/home/sysadmin/.hermes/logs/voice-agent-stdout.log
StandardError=append:/home/sysadmin/.hermes/logs/voice-agent-stderr.log

[Install]
WantedBy=multi-user.target
